# rootfs/etc/services.d/refresh/run
#!/usr/bin/with-contenv bash
echo "Starting refresh process."
exec /opt/scripts/refresh.sh

